# crontab
0 0 * * * poetry run python /app/etl/etl.py