# crontab
0 0 * * * /app/.venv/bin/python /app/etl/run_daily_etl.py